# karto/config.dist/cron/crontabs/root
# min   hour    day     month   weekday command
*/15    *       *       *       *       busybox run-parts ${CRONPERIODIC}/15min
5       *       *       *       *       busybox run-parts ${CRONPERIODIC}/hourly
10      2       *       *       *       busybox run-parts ${CRONPERIODIC}/daily
10      3       *       *       6       busybox run-parts ${CRONPERIODIC}/weekly
10      5       1       *       *       busybox run-parts ${CRONPERIODIC}/monthly

